# auto/auto/backup-email-cloudfiles.timer
[Unit]
Description=Backup email to cloudfiles

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target