# auto/auto/getmail.timer
[Unit]
Description=Downloads email

[Timer]
OnBootSec=60sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target