# auto/auto/restart-network.timer
[Unit]
Description=Disconnects and reconnects the ethernet connection

[Timer]
OnBootSec=23h
OnUnitActiveSec=23h

[Install]
WantedBy=timers.target